# setup/configs/supervisor/supervisord-view.conf
[program:view]
command=/start-view.sh
numprocs=1
autostart=true
autorestart=true