# setup/configs/supervisor/supervisord-camel.conf
[program:camel]
command=/start-camel.sh
numprocs=1
autostart=true
autorestart=true